# tommy/rootfs/etc/services.d/tommy/run
#!/usr/bin/with-contenv bashio

# Get the ingress URL from Supervisor API
INGRESS_URL=$(bashio::addon.ingress_url || echo "")
bashio::log.info "Ingress URL: ${INGRESS_URL}"

# Extract just the path part if we got a URL
if [ -n "$INGRESS_URL" ]; then
    INGRESS_PATH=$(echo "$INGRESS_URL" | sed 's|^https\?://[^/]*||')
    INGRESS_PATH="${INGRESS_PATH%/}"
    export INGRESS_PATH
    bashio::log.info "Extracted ingress path: ${INGRESS_PATH}"
fi

# Start virtual bridge
if [ -n "$INGRESS_PATH" ]; then
    cd /app || exit
    exec ./start.sh --ingress-path "$INGRESS_PATH"
else
    cd /app || exit
    exec ./start.sh
fi